# files/etc/supervisor.d/subversion-edge.conf
[program:csvn]
command=/opt/csvn/bin/csvn console 

[program:csvn-httpd]
command=/opt/csvn/bin/csvn-httpd start
